# github-dotfile.service
[Unit]
Description=GitHub dotfile service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/github-dotfiles

[Install]
WantedBy=default.target
